# packages/grid/seaweedfs/config/app/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:seaweed_server]
priority=1
command=sh -c "exec weed server \
    -idleTimeout=60 \
    -master -master.dir=$SWFS_MASTER_DIR -master.volumeSizeLimitMB=$SWFS_VOLUME_SIZE_LIMIT_MB \
    -filer -filer.maxMB=$SWFS_FILER_CHUNKS_MB -filer.concurrentUploadLimitMB=$SWFS_FILER_UPLOAD_LIMIT_MB \
    -s3 -s3.config=/run/secrets/seaweedfs/s3config.json \
    -volume -dir=$SWFS_VOLUME_DIR -volume.dir.idx=$SWFS_VOLUME_IDX_DIR -volume.max=$SWFS_VOLUME_MAX -volume.index=$SWFS_VOLUME_INDEX \
    -volume.concurrentUploadLimitMB=$SWFS_VOLUME_UPLOAD_LIMIT_MB -volume.concurrentDownloadLimitMB=$SWFS_VOLUME_DOWNLOAD_LIMIT_MB"
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:mount_api_server]
priority=2
command=sh -c ". scripts/wait_for_swfs.sh && exec uvicorn --host=0.0.0.0 --port=$MOUNT_API_PORT --log-config=logging.yaml src.api:app"
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[include]
files = /data/mounts/**/*.conf /etc/mounts/**/*.conf
