# packages/grid/seaweedfs/config/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:seaweed_server]
command=bash -c "weed server \
    -idleTimeout=60 \
    -master -master.dir=$SWFS_MASTER_DIR -master.volumeSizeLimitMB=$SWFS_VOLUME_SIZE_LIMIT_MB \
    -filer -filer.maxMB=64 -filer.concurrentUploadLimitMB=0 \
    -s3 -s3.config=$S3_CONFIG \
    -volume -dir=$SWFS_VOLUME_DIR -volume.dir.idx=$SWFS_VOLUME_IDX_DIR -volume.max=$SWFS_VOLUME_MAX -volume.index=leveldb \
    -volume.concurrentUploadLimitMB=0 -volume.concurrentDownloadLimitMB=0"
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:mount_api_server]
command=uvicorn --host=0.0.0.0 --port=%(ENV_MOUNT_API_PORT)s --log-level=%(ENV_UVICORN_LOG_LEVEL)s src.api:app
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[include]
files = /data/mount/procs/*.conf ~/mount/*.conf
