# packages/grid/seaweedfs/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:seaweed_master]
command=weed master -mdir=/data/master -volumeSizeLimitMB=%(ENV_S3_VOLUME_SIZE_MB)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:seaweed_filer]
command=weed filer -master=0.0.0.0:9333 -s3 -s3.port=8333 -s3.config=%(ENV_S3_CONFIG_PATH)s -ip=0.0.0.0 -port=8888 -dataCenter=dc1 -rack=rack1
autorestart=unexpected
exitcodes=011
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:seaweed_volume]
command=weed volume -dir=%(ENV_S3_MOUNT_DIRS)s -mserver=0.0.0.0:9333 -max=%(ENV_S3_MOUNT_DIRS_MAX)s -dataCenter=dc1 -rack=rack1
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

# [program:seaweed_server]
# command=bash -c "weed server -master.mdir=/data/master -dir=/data/vol0 -s3 -s3.port=$S3_PORT -volume.max=500 -master.volumeSizeLimitMB=$S3_VOLUME_SIZE_MB"
# autorestart=unexpected
# exitcodes=0
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# redirect_stderr=true

# [program:configure_s3]
# command=bash -c "sleep 15; echo s3.configure -access_key=$S3_ROOT_USER -secret_key=$S3_ROOT_PWD -user=iam -actions=Read,Write,List,Tagging,Admin -apply | weed shell && echo Root IAM Configured."
# autorestart=unexpected
# exitcodes=0
# stdout_logfile=/tmp/configure_s3.log
# stdout_logfile_maxbytes=0
# redirect_stderr=true

[program:mount_api_server]
command=uvicorn --host=0.0.0.0 --port=%(ENV_SEAWEED_MOUNT_PORT)s --log-level=%(ENV_UVICORN_LOG_LEVEL)s src.app:app
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[include]
files = /data/procs/*.conf
